# sma-bluetooth/src/smareader.timer
[Unit]
Description=Run smareader.service every 60 seconds 
Requires=smareader.service

[Timer]
OnCalendar=*:*:00

[Install]
WantedBy=timers.target
